# dkms_autoinstaller.upstart
# dkms_autoinstaller - Builds all DKMS modules for current kernel if necessary
#
# Run DKMS prior to GDM/KDM/Friends

description	"DKMS Automatic Build & Installation Service"
author          "Mario Limonciello <Mario_Limonciello@Dell.com>"

start on (startup and filesystem)

task

emits build-successful
emits build-failed

exec /usr/lib/dkms/dkms_autoinstaller start $KERNEL
