# dkms_autoinstaller.upstart
# dkms_autoinstaller - Builds all DKMS modules for current kernel if necessary
#
# Run DKMS prior to GDM/KDM/Friends

description	"DKMS Automatic Build & Installation Service"
author          "Mario Limonciello <Mario_Limonciello@Dell.com>"

start on (starting gdm
          or starting kdm
          or starting xdm
          or starting oem-config
          or starting ubiquity
          or stopping rc RUNLEVEL=[2345])

task

exec /usr/lib/dkms/dkms_autoinstaller start
